# deb/prerm
#!/bin/bash

echo "prerm"

systemctl disable neuron
systemctl stop neuron

rm -rf /opt/neuron/core/*
rm -rf /opt/neuron/logs/*
